# docker/php-worker/conf.d/php-worker.conf.example
[program:php-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/bin/php-worker.php
autostart=true
autorestart=true
numprocs=1
